# scripts/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:frps]
command=/opt/frp/frps -c /etc/frp/frps.toml
user=frp
group=frp
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/frp/frps.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
environment=HOME="/opt/frp",USER="frp"

[program:cron]
command=/usr/sbin/cron -f
user=root
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/cron.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=2

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface
